# deploy/altin-report.timer
[Unit]
Description=Altin gun sonu raporu (her gun 18:45 TR)

[Timer]
# TR = UTC+3. 18:45 TR = 15:45 UTC. EVDS (15:30) bittikten sonra calisir.
OnCalendar=*-*-* 15:45:00
Persistent=true

[Install]
WantedBy=timers.target
