# deploy/altin-evds.service
[Unit]
Description=Altin Takip - Gunluk EVDS guncelleme
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/altin
ExecStart=/home/ubuntu/altin/.venv/bin/python -m src.evds_job daily
