# deploy/altin-bot.service
[Unit]
Description=Altin Takip - Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/altin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/altin/.venv/bin/python -m src.telegram_bot
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
