# deploy/altin-backup.timer
[Unit]
Description=Altin gunluk yedek (opsiyonel; her gun 02:30 TR)

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true

[Install]
WantedBy=timers.target
