# deploy/altin-collector.service
[Unit]
Description=Altin Takip - Veri Toplayici
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/altin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/altin/.venv/bin/python -m src.collector
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
